# aur-update-npm-package.timer
[Unit]
Description=Run update-aur-npm-package weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
